# src/pybot/youpi2/setup/pkg_data/youpi2.target
[Unit]
Description=Youpi2 sub-system

[Install]
WantedBy=multi-user.target
